# .claude/skills/cron-register/template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.telegram-skill-bot.telegram-agent.__NAME__</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__ROOT__/routines/push.sh</string>
        <string>--model</string>
        <string>haiku</string>
        <string>--prompt</string>
        <string>__PROMPT__ (톤규칙 필수: 사용자 호칭, 존댓말 반말금지, ** 안씀, 기호 최소. 메시지 본문만 출력)</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>__HOUR__</integer>
        <key>Minute</key>
        <integer>__MINUTE__</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__ROOT__/runtime/logs/__NAME__.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/runtime/logs/__NAME__.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>__ROOT__</string>
</dict>
</plist>
