# agents/.template/routines/com.telegram-skill-bot.telegram-agent.mirror-poll.timer
[Unit]
# DGN-268 S3: fires the paired mirror-poll.service every 300s (5-min cadence,
# matching the macOS plist StartInterval 300). Not RunAtLoad -- OnUnitActiveSec
# schedules the FIRST run 300s after the timer starts, same as the plist
# (RunAtLoad false). Persistent runs a missed cycle after sleep/downtime.
Description=dogany mirror poll timer (__AGENT_NAME__) -- 300s cadence

[Timer]
OnUnitActiveSec=300
OnBootSec=300
Persistent=true
Unit=com.telegram-skill-bot.__AGENT_NAME__.mirror-poll.service

[Install]
WantedBy=timers.target
