# agents/.template/routines/com.telegram-skill-bot.telegram-agent.cleanup-files.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.telegram-skill-bot.__AGENT_NAME__.cleanup-files</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_ROOT__/routines/cleanup-files.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/.telegram_bot/logs/cleanup-files.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/.telegram_bot/logs/cleanup-files.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>
</dict>
</plist>
